# .config/systemd/user/udiskie.service
[Unit]
Description=Automounting Daemon

[Service]
ExecStart=/usr/bin/udiskie -qN

[Install]
WantedBy=default.target
