# build/intelbacklight/99-intelbacklight.rules
# Adjust screen brightness according to power state
# On AC
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ATTR{online}=="1", RUN+="/usr/bin/intelbacklight.sh 90"
# On Battery
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ATTR{online}=="0", RUN+="/usr/bin/intelbacklight.sh 20"
